# pipeline/crontab
SHELL=/bin/bash

# Heartbeat
* * * * * date '+[heartbeat] %Y-%m-%d %H:%M:%S'

# Pipeline schedules
*/5 * * * * flock -n /tmp/auto_merge.lock sh -c 'cd /app && python auto_merge.py'
0 2 * * *   sh -c 'cd /app && python auto_cleaning_and_dq.py'
15 2 * * *  sh -c 'cd /app && python auto_dq_report.py'